# images/cron/crontab
*/5 * * * * /bin/bash /root/top10.sh
# Single container cron example
